# ai-assistant.service
[Unit]
Description=AI Personal Assistant Service
After=network.target

[Service]
Type=simple
User=tyson
WorkingDirectory=/home/tyson/AiTools/website-monitoring-tool
ExecStart=/home/tyson/miniconda3/bin/python3 /home/tyson/AiTools/website-monitoring-tool/monitor.py
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
